# rabbit/supervisord.conf
[supervisord]
nodaemon = true
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
startsecs=3
[program:rabbitmq-server]
command=/usr/sbin/rabbitmq-server
autostart=true
startsecs=30
