# rabbit/supervisord.conf
[supervisord]
nodaemon = true
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
startsecs=3
[program:glance-api]
command=/usr/sbin/rabbitmq-server
autostart=true
user=rabbitmq
directory=/var/lib/rabbitmq
startsecs=30
