# keystone/supervisord.conf
[supervisord]
nodaemon = true
[program:sshd]
command=sshd -D
autostart=true
startsecs=3
[program:keystone]
command=keystone-all
autostart=true
user=keystone
directory=/var/lib/keystone
startsecs=30

