# keystone/supervisord.conf
[supervisord]
nodaemon = true
[program:keystone]
command=keystone-all
autostart=true
user=keystone
directory=/var/lib/keystone
startsecs=30

